# etc/profile.d/antigen.sh
#!/usr/bin/env bash

if [ -f "/usr/share/zsh/share/antigen.zsh" ] && [ "$0" == "zsh" ]; then
    
    source "/usr/share/zsh/share/antigen.zsh"

fi

