# etc/profile.d/direnv.sh
#!/usr/bin/env bash

if [ "$(command -v direnv)" ]; then
    
    eval "$(direnv hook bash)"
    
fi
